# uwsgi.upstart.conf
description "uWSGI server configured to serve theatertilt"

start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/local/bin/uwsgi --ini /etc/uwsgi/theatertilt.uwsgi.ini
